# packages/net-misc/connman/files/connman.initd
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Purpose License v2
# $Header: $

depend() {
	need dbus localmount
	after bootmisc modules
}

start() {
	ebegin "Starting Connection Manager daemon"
	start-stop-daemon --start --exec /usr/sbin/connmand
	eend $?
}

stop() {
	ebegin "Stopping Connection Manager daemon"
	start-stop-daemon --stop --exec /usr/sbin/connmand
	eend $?
}
